# supervisord/deluged.conf
[program:deluged]
command=/usr/bin/deluged
