# supervisord/flexget.conf
[program:flexget]
command=/usr/local/bin/flexget
